# src/etc/rc.d/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

find /etc/rc.d/rc.local.d/ -type f -exec /bin/sh {} \;

touch /var/lock/subsys/local
